# images/proxy/services.d/httpd/run
#!/usr/bin/with-contenv bash
/usr/sbin/httpd -DFOREGROUND